# supervisord.conf
[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=debug               ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=true              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket
            ; override serverurl computation (childutils)

[group:tornadoes]
programs = tornado-8000, tornado-8001, tornado-8002, tornado-8003, tornado-8004, tornado-8005

[program:tornado-8000]
command=python3 app.py --port=8000
directory=/home/solomon/Desktop/tornado/
user=solomon
autorestart=true
redirect_stderr=true
stdout_logfile=/home/solomon/tornado.log
loglevel=info
stopsignal=QUIT

[program:tornado-8001]
command=python3 app.py --port=8001
directory=/home/solomon/Desktop/tornado/
user=solomon
autorestart=true
redirect_stderr=true
stdout_logfile=/home/solomon/tornado.log
loglevel=info
stopsignal=QUIT 

[program:tornado-8002]
command=python3 app.py --port=8002
directory=/home/solomon/Desktop/tornado/
user=solomon
autorestart=true
redirect_stderr=true
stdout_logfile=/home/solomon/tornado.log
loglevel=info
stopsignal=QUIT 

[program:tornado-8003]
command=python3 app.py --port=8003
directory=/home/solomon/Desktop/tornado/
user=solomon
autorestart=true
redirect_stderr=true
stdout_logfile=/home/solomon/tornado.log
loglevel=info
stopsignal=QUIT 

[program:tornado-8004]
command=python3 app.py --port=8004
directory=/home/solomon/Desktop/tornado/
user=solomon
autorestart=true
redirect_stderr=true
stdout_logfile=/home/solomon/tornado.log
loglevel=info
stopsignal=QUIT 

[program:tornado-8005]
command=python3 app.py --port=8005
directory=/home/solomon/Desktop/tornado/
user=solomon
autorestart=true
redirect_stderr=true
stdout_logfile=/home/solomon/tornado.log
loglevel=info
stopsignal=QUIT 

